# dev/services/certbot.timer
[Unit]
Description=Check SSL certs every minute

[Timer]
OnActiveSec=60

